# bb-wlan0-defaults/suite/bullseye/debian/81-add-SoftAp0-interface.rules
SUBSYSTEM=="net", ACTION=="add", KERNEL=="mlan[0-9]", DRIVERS=="mwifiex_pcie" \
	RUN+="/sbin/iw dev %k interface add SoftAp0 type __ap"
