# generic-sys-mods/suite/noble/debian/generic-sys-mods.beagle-camera-setup.service
[Unit]
Description=BeagleBoard Camera Setup
ConditionFileIsExecutable=/usr/sbin/beagle-camera-setup

[Service]
Type=oneshot
ExecStart=/usr/sbin/beagle-camera-setup

[Install]
WantedBy=multi-user.target
